# gemba-validator/auto/systemd/gemba-auto-unjail.service
[Unit]
Description=GembaBlockchain validator watchdog: detect-stuck -> restart -> sync -> unjail (one-shot, timer-driven)
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/gemba/validator-auto.env
ExecStart=/usr/local/bin/gemba-auto-unjail.sh
Nice=10
